# telegram-terminal-bot.service
[Unit]
Description=Telegram Terminal Bot
After=network.target

[Service]
Type=simple
User=kegon
WorkingDirectory=/home/kegon/prog/secure-telegram-terminal
ExecStartPre=/usr/bin/mkdir -p /home/kegon/prog/secure-telegram-terminal/logs /home/kegon/prog/secure-telegram-terminal/data
ExecStart=/home/kegon/prog/secure-telegram-terminal/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
