# telegram-terminal-bot.service
[Unit]
Description=Telegram Terminal Bot
After=network.target

[Service]
Type=simple
User=kegon
WorkingDirectory=/mnt/d/bot_telegram_trminal_logs
Environment=PATH=/mnt/d/bot_telegram_trminal_logs/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -c 'cd /mnt/d/bot_telegram_trminal_logs && python3 -m venv venv'
ExecStartPre=/bin/bash -c 'cd /mnt/d/bot_telegram_trminal_logs && ./venv/bin/pip install -q -r requirements.txt'
ExecStartPre=/bin/bash -c 'mkdir -p /mnt/d/bot_telegram_trminal_logs/logs'
ExecStartPre=/bin/bash -c 'tmux has-session -t claude 2>/dev/null || SHELL=/bin/bash tmux new-session -d -s claude -c /mnt/d/bot_telegram_trminal_logs'
ExecStart=/bin/bash -c 'cd /mnt/d/bot_telegram_trminal_logs && source venv/bin/activate && python bot.py'
Restart=always
RestartSec=10
StandardOutput=append:/mnt/d/bot_telegram_trminal_logs/logs/bot.log
StandardError=append:/mnt/d/bot_telegram_trminal_logs/logs/bot_error.log

[Install]
WantedBy=multi-user.target
